# rc.local
#!/bin/bash
cp /var/setup_ss/quick_ss/mincn.txt /dev/shm
ss_watch_dog
exit 0
